# autostart/piframe.desktop
[Desktop Entry]
Type=Application
Name=piframe
Exec=/home/peter/piframe/start_browser.sh
